# mlproject/supervisord.conf
[supervisord]
nodaemon=true

[program:mlflow]
command=mlflow ui --host 0.0.0.0 --port 9610 --backend-store-uri %(ENV_MLFLOW_TRACKING_URI)s --default-artifact-root file:///mlflow/artifacts
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:app]
command=python3 app.py
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0